# crontab
0   22  *   *   *   "$HOME/.scripts/backup" &> "$HOME/backup.txt" && "$HOME/.scripts/i3status_refresh"
0   22  *   *   *   checkupdates | wc -l > "$HOME/.updates"
