# crontab
0   22  *   *   *   /home/vuk/.scripts/backup.sh &> /home/vuk/backup.txt
